# configurations/docker/crontab/entrypoint.sh
#!/bin/bash

crontab producer.cron
printenv >> /etc/environment
cron -f
